# packages/debian/DEBIAN/prerm
#!/bin/bash

set -e

deb-systemd-helper stop cardwired.service

deb-systemd-helper disable cardwired.service